# systemd/maskforai.service.d/10-adaptive-proxy.conf
[Service]
RuntimeDirectory=maskforai
RuntimeDirectoryMode=0700
Environment=MASKFORAI_ADAPTIVE_PROXY=1
ExecStartPre=__MASKFORAI_DETECTOR_PATH__ %t/maskforai/proxy.env
EnvironmentFile=-%t/maskforai/proxy.env
